# provisioning/common/services/cozy-presence.service
[Unit]
Description=cozy-presence: Local identity persistence daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User={{ pillar.get('users', {}).get(managed_users[0], {}).get('name', 'vegcom') }}
Environment=PYTHONPATH=/opt/cozy-presence/src
ExecStart=/opt/miniconda/envs/cozy-presence/bin/python -c "
import sys, time
sys.path.insert(0, '/opt/cozy-presence/src')
from cozy_presence.store import Store
from cozy_presence.index import Index
from cozy_presence.sync import SyncClient

# Initialize components
store = Store()
index = Index()
sync = SyncClient()

print('cozy-presence daemon started')

# Bootstrap index on startup
if index.available:
    entities = store.all_entities()
    if entities:
        count = index.bootstrap(entities)
        print(f'Bootstrapped {count} entities')

# Keep alive
while True:
    time.sleep(60)
    # TODO: periodic sync, health checks, etc.
"
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target