# provisioning/linux/files/etc-skel/.config/systemd/user/input-leapd.service
# ~/.config/systemd/user/input-leapd.service
[Unit]
Description=Leap Input Server
Documentation=https://github.com/input-leap/input-leap
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/input-leap --daemon
Restart=on-failure
RestartSec=2
# Optional: ensure XDG paths are respected
Environment=XDG_CONFIG_HOME=%h/.config
Environment=XDG_DATA_HOME=%h/.local/share

[Install]
WantedBy=default.target
Requires=graphical.target
