# provisioning/linux/files/etc-skel/.config/systemd/user/plasma-session.target
[Unit]
Description=Custom Return-to-Desktop handler
Requires=return-to-desktop.service
After=return-to-desktop.service
