# provisioning/linux/files/etc-skel/.config/systemd/user/return-to-desktop.service
[Unit]
Description=Return to Desktop override

[Service]
Type=oneshot
ExecStart=/usr/local/bin/return-to-desktop
