# provisioning/linux/files/etc-systemd-user/cozy-presence@.service
[Unit]
Description=cozy-presence: Local identity persistence daemon
After=network.target
Wants=network.target

[Service]
Type=simple
Environment=PYTHONPATH=/opt/cozy/cozy-presence/src
ExecStart=/opt/miniforge3/envs/cozy-presence/bin/python /opt/cozy/cozy-presence/bin/cozy-presence-daemon.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
